# installer/macos/com.studiolink.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.studiolink.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/studiolink-daemon</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/StudioLink/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/StudioLink/daemon.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PLUGIN_PORT</key>
    <string>45678</string>
  </dict>
</dict>
</plist>
